# backend/shadowcraft-engine-all.conf
start on stopped cloud-init
stop on shutdown

env NUM_WORKERS=4

pre-start script
  for i in `seq 1 $NUM_WORKERS`
  do
    start shadowcraft-engine N=$i
  done
end script
